# supplements/ansible/roles/rtb_latency/templates/etc/systemd/system/rtb_latency.service.j2
[Unit]
Description=RTB Latency Checker
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python2 {{lc_install_dir}}/rtb_latency.py

[Install]
WantedBy=network-online.target
