# supplements/systemd/latency_check.timer
[Unit]
Description=Latency Check

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
