# supplements/cron/latency_check.cron
* * * * * /usr/bin/python2 /opt/latency_check/latency_check.py | /usr/bin/logger -s
